# tests/thinking/node-1/supervisord.conf
; Supervisord Configuration file

[inet_http_server]
port=*:%(ENV_SUPERVISORD_HTTP_PORT)s
username=%(ENV_SUPERVISORD_HTTP_USER)s
password=%(ENV_SUPERVISORD_HTTP_PASSWORD)s

[supervisord]
logfile=/toskose/supervisord/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=%(ENV_SUPERVISORD_LOG_LEVEL)s
pidfile=/toskose/supervisord/tmp/supervisord.pid
nodaemon=true
minfds=1024
minprocs=200
umask=022
user=root
identifier=supervisor
nocleanup=true
childlogdir=/toskose/supervisord/tmp
strip_ansi=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; api configure
[program:api-configure]
command=/bin/sh -c '/toskose/apps/api/scripts/configure.sh'
process_name=api-configure
numprocs=1
umask=022
priority=999
autostart=false
startsecs=0
startretries=3
autorestart=false
exitcodes=0
stopsignal=TERM
stopwaitsecs=10
stopasgroup=false
killasgroup=false
user=root
redirect_stderr=true
stdout_logfile=/toskose/apps/api/logs/api-configure.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
serverurl=AUTO

; api install
[program:api-install]
command=/bin/sh -c '/toskose/apps/api/scripts/install.sh'
process_name=api-install
numprocs=1
umask=022
priority=999
autostart=false
startsecs=0
startretries=3
autorestart=false
exitcodes=0
stopsignal=TERM
stopwaitsecs=10
stopasgroup=false
killasgroup=false
user=root
redirect_stderr=true
stdout_logfile=/toskose/apps/api/logs/api-install.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
serverurl=AUTO

; api start
[program:api-start]
command=/bin/sh -c '/toskose/apps/api/scripts/start.sh'
process_name=api-start
numprocs=1
umask=022
priority=999
autostart=false
startsecs=0
startretries=3
autorestart=false
exitcodes=0
stopsignal=TERM
stopwaitsecs=10
stopasgroup=false
killasgroup=false
user=root
redirect_stderr=true
stdout_logfile=/toskose/apps/api/logs/test-1.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
serverurl=AUTO

; api stop
[program:api-stop]
command=/bin/sh -c '/toskose/apps/api/scripts/stop.sh'
process_name=api-stop
numprocs=1
umask=022
priority=999
autostart=false
startsecs=0
startretries=3
autorestart=false
exitcodes=0
stopsignal=TERM
stopwaitsecs=10
stopasgroup=false
killasgroup=false
user=root
redirect_stderr=true
stdout_logfile=/toskose/apps/api/logs/api-stop.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
serverurl=AUTO

; api uninstall
[program:api-uninstall]
command=/bin/sh -c '/toskose/apps/api/scripts/uninstall.sh'
process_name=api-uninstall
numprocs=1
umask=022
priority=999
autostart=false
startsecs=0
startretries=3
autorestart=false
exitcodes=0
stopsignal=TERM
stopwaitsecs=10
stopasgroup=false
killasgroup=false
user=root
redirect_stderr=true
stdout_logfile=/toskose/apps/api/logs/api-uninstall.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
serverurl=AUTO
